# contrib/systemd/mcp-proxy.service
[Unit]
Description=MCP Proxy - HTTP/SSE gateway for MCP servers
After=network.target

[Service]
Type=simple
ExecStart=/opt/mcp-proxy/bin/mcp-proxy \
    --host 0.0.0.0 \
    --port 8000 \
    --allow-origin "*" \
    --named-server-config /etc/ozwell/mcp/servers.json
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
