# contrib/systemd/openvscode-server.service
[Unit]
Description=OpenVSCode Server - VS Code in the browser
After=network.target continue-env.service

[Service]
Type=simple
WorkingDirectory=/workspace
Environment=SHELL=/bin/bash
Environment=HOME=/root
Environment=USER=root
Environment=TERM=xterm-256color
PassEnvironment=ANTHROPIC_API_KEY OPENAI_API_KEY
ExecStart=/opt/openvscode-server/bin/openvscode-server --host 127.0.0.1 --port 8081 --server-base-path /vscode --without-connection-token --telemetry-level off
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
