# contrib/systemd/ttyd.service
[Unit]
Description=ttyd - Web terminal
After=network.target

[Service]
Type=simple
WorkingDirectory=/workspace
ExecStart=/usr/local/bin/ttyd -W -p 7682 -i lo -t disableLeaveAlert=true tmux new-session -A -s main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
